# debian/postrm
#!/bin/sh

#DEBHELPER#

set -e

rm -f /etc/nginx/sites-enabled/debian-timeline
invoke-rc.d nginx reload
